# cookbooks/ntpdate/files/default/etc/cron.d/ntpdate
0 */3 * * *     root    /usr/sbin/ntpdate 2.rs.pool.ntp.org
